# services/systemd/fw-maker-bot.service
[Unit]
Description=Freimarkets test maker-bot (always-online P2P maker, own keys, FRC->BTC)
After=network.target fw-market-relay.service

[Service]
Type=simple
WorkingDirectory=/root/free-money/freicoin-wallet
ExecStart=/usr/bin/node /root/nv3-playground/maker-bot.mjs
Restart=always
RestartSec=5
StandardOutput=append:/root/nv3-playground/maker-bot.log
StandardError=append:/root/nv3-playground/maker-bot.log

[Install]
WantedBy=multi-user.target
